# gnome-search-provider/drillbit-search-provider.service
[Unit]
Description=Drillbit GNOME Shell Search Provider
After=graphical-session.target

[Service]
Type=dbus
BusName=org.drillbit.SearchProvider
ExecStart=SCRIPT_PATH
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
